# deploy/config/supervisord.ini
[program:messenger-consume]
command=php /home/huy/api/current/apps/api/bin/console messenger:consume async --time-limit=300
user=ubuntu
numprocs=2
startsecs=0
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
